# herdr/.config/systemd/user/herdr.service
[Unit]
Description=Herdr persistent agent multiplexer

[Service]
Type=simple
ExecStart=%h/.local/bin/herdr server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
